# deploy/systemd/mpcontrol-backup.service
[Unit]
Description=MP Control daily backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/mpcontrol
ExecStart=/bin/bash /opt/mpcontrol/scripts/backup_daily.sh
EnvironmentFile=-/opt/mpcontrol/.env
User=root
StandardOutput=journal
StandardError=journal
